# rootfs/etc/services.d/ha_inetbox/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Runs truma_service
# ==============================================================================

# Wait mosquitto
bashio::net.wait_for $(bashio::config 'port')

# Start the inetbox service
bashio::log.info "Starting Truma Inetbox service..."
exec cd /etc/services.d/ && python -m  ha_inetbox.homeassistant_truma_service